# scripts/mcp-agent.service
# systemd service template for mcp-agent.
# install.sh substitutes <install-user> and <AGENT_HOME> before writing this to
# /etc/systemd/system/mcp-agent.service.
[Unit]
Description=MCP Agent server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=<install-user>
Environment=AGENT_HOME=<AGENT_HOME>
EnvironmentFile=<AGENT_HOME>/agent.env
# Prepend the venv so the web-search MCP subprocess finds python3 + mcp.
Environment=PATH=<AGENT_HOME>/venv/bin:/usr/local/bin:/usr/bin:/bin
WorkingDirectory=<AGENT_HOME>
ExecStart=/usr/bin/java -jar <AGENT_HOME>/bin/mcp-integration.jar \
  --spring.config.additional-location=file:<AGENT_HOME>/config/application.yml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
